# templates/codeshield-secrets.service
[Unit]
Description=CODE SHIELD V3 — Decrypt secrets to tmpfs
Documentation=https://github.com/godlovestome/codeshield_claude
Before=openclaw.service
Before=openclaw-watcher.service

[Service]
Type=oneshot
RemainAfterExit=yes

# RuntimeDirectory creates /run/openclaw-codeshield (tmpfs-backed)
# Automatically cleaned up when service stops
RuntimeDirectory=openclaw-codeshield
RuntimeDirectoryMode=0700
RuntimeDirectoryPreserve=yes

ExecStart=/usr/local/sbin/codeshield-secrets-unseal
ExecStop=/bin/rm -f /run/openclaw-codeshield/secrets.env

StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
